# consul-base/supervisord.conf
[supervisord]
nodaemon=true


[program:consul]
command=/usr/bin/consul agent -config-dir=/etc/consul.d/






